# production_settings/gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket - Tanks

[Socket]
ListenStream=unix:/home/dimtasev/tanks.sock

[Install]
WantedBy=sockets.target
